# system/etc/supervisor/conf.d/ei.conf
[program:ei]
command=/home/pi/ecstatic-epiphany/ei
autostart=true
autorestart=true
redirect_stderr=true
directory=/home/pi/ecstatic-epiphany
user=pi
