# system/etc/supervisor/conf.d/tunnel.conf
[program:ei]
command=/usr/bin/ssh ei@54.196.4.174 -R 11586:localhost:22 -N
autostart=true
autorestart=true
redirect_stderr=true
user=micah

