# system/etc/supervisor/conf.d/tunnel.conf
[program:tunnel]
command=/usr/bin/ssh ei@pixelbees.net -R 0:localhost:22 -N
autostart=true
autorestart=true
redirect_stderr=true
user=micah

